# etc/supervisor/conf.d/exampled.conf
[program:exampled]
command = /php-project/exampled.initd.sh.php
user = root
autostart = false
